# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm7.1]
command=/usr/local/php/sbin/php-fpm
#command=/usr/local/php/sbin/php-fpm  -c /usr/local/php/etc/php-fpm.conf
user=www-data
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/php-fpm.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/php-fpm-error.log
stderr_logfile_maxbytes=10MB

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true

